# crontab
@reboot /usr/bin/screen -c /root/screen-printermonitor -d -m
